# PAI/PULSE/linux/pai-cadence-weekly.timer
[Unit]
Description=PAI weekly cadence (Mon 07:00)

[Timer]
OnCalendar=Mon 07:00
Persistent=true

[Install]
WantedBy=timers.target
